# runatstartup.conf
description     "Executes startup script at Rocket Repo"
start on startup
task
exec /home/pi/Rocket/shstart.sh
